# deploy/system.service
[Unit]
Description=La Burbuja POS (Next.js)
After=network-online.target nodered.service
Wants=network-online.target nodered.service

[Service]
Type=simple
User=mdares
Group=mdares
WorkingDirectory=/home/mdares/Laundry
Environment=NODE_ENV=production
Environment=PORT=3000
EnvironmentFile=/home/mdares/Laundry/.env
ExecStartPre=/usr/bin/npm run prisma:migrate:deploy
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
TimeoutStartSec=90
KillSignal=SIGINT
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
